# share/etc/xdg/autostart/setgnomekbd.desktop
[Desktop Entry]
Name=Set keyboard
Type=Application
Exec=setgnomekeyboard
Terminal=false
Categories=System;
